# app.conf

[supervisord]
nodaemon=true

[program:cron]
command=/usr/sbin/cron -f
stdout_events_enabled=true
stderr_events_enabled=true
stdout_logfile=NONE
stderr_logfile=NONE

[program:docker-gen]
command=/usr/local/bin/docker-gen -notify "/usr/bin/supervisorctl start logrotate-debug" -watch /app/logrotate.tmpl /etc/logrotate.d/docker
stdout_events_enabled=true
stderr_events_enabled=true
stdout_logfile=NONE
stderr_logfile=NONE

[program:logrotate-debug]
command=/usr/sbin/logrotate /etc/logrotate.conf -d
autorestart=false
startsecs=0
startretries=0
stdout_events_enabled=true
stderr_events_enabled=true
stdout_logfile=NONE
stderr_logfile=NONE


[eventlistener:stdout]
command = supervisor_stdout
buffer_size = 100
events = PROCESS_LOG
result_handler = supervisor_stdout:event_handler
